# run.ini
[program:app]
command=/usr/bin/node /opt/hook/app.js
environment=NODE_ENV=production
directory=/opt/hook/
autostart=true
autorestart=true
